# deploy/supervisor-cron-app.ini
[program:dkron]
command=/data/apps/dkron/bin/dkron agent --config /data/apps/dkron/dkron.json --mail-password=%(ENV_MAIL_PASSWD)s
user=root
autostart=true
autorestart=true
stdout_logfile=/qke/log/dkron_stdout.log
stderr_logfile=/qke/log/dkron_stderr.log
stopasgroup=true
stopsignal=QUIT